# infra/systemd/studybox-backup.timer
[Unit]
Description=Run StudyBox AI PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=studybox-backup.service

[Install]
WantedBy=timers.target
